# templates/dcgmd-telegraf.service.epp
<%- | $exec_start,
      $dcgm_telegraf_py_path,
| -%>
# This file is managed by Puppet, manual changes will be lost
[Unit]
Description=DCGM Telegraf service
Requires=nvidia-dcgm.service
 
[Service]
Type=simple
User=telegraf
PrivateTmp=false
ExecStart=<%= $exec_start %> <%= $dcgm_telegraf_py_path %> -i 60
Restart=on-abort
 
[Install]
WantedBy=multi-user.target

