# files/dcgmd-telegraf.service.epp
<%- | $execStart,
| -%>
# This file is managed by Puppet, manual changes will be lost
[Unit]
Description=DCGM Telegraf service
Requires=nvidia-dcgm.service
 
[Service]
Type=simple
User=telegraf
PrivateTmp=false
ExecStart=<%= $execStart %>
Restart=on-abort
 
[Install]
WantedBy=multi-user.target

